# restic@.timer
[Unit]
Description=Run Restic hourly

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
