# systemd/restic-prune@.timer
[Unit]
Description=Run Restic prune on %I weekly

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
