# systemd/restic-prune@.timer
[Unit]
Description=Run Restic prune on %I weekly

[Timer]
# OnCalendar=*-*-* 0:00:00
OnCalendar=weekly

[Install]
WantedBy=timers.target
